# icinga2.service.d/10-loglevel.conf
[Service]
# need to unset ExecStart once because it would try to array it.
ExecStart=
ExecStart=/usr/sbin/icinga2 daemon -e /var/log/icinga2/icinga2.err -x warning
